# build/docker/supervisord.conf
[supervisord]
nodaemon=true
silent=false
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:node]
command=node server.js
directory=%(ENV_FALKORDB_BROWSER_PATH)s
user=nextjs
autostart=%(ENV_BROWSER)s
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/nextjs",USER="nextjs"

[program:text-to-cypher]
command=bash -c 'set -o pipefail; ./text-to-cypher 2>&1 | while IFS= read -r line; do echo "[$(date "+%%Y-%%m-%%d %%H:%%M:%%S")] [TEXT-TO-CYPHER] $line"; done'
directory=%(ENV_FALKORDB_BROWSER_PATH)s
autostart=%(ENV_CYPHER)s
user=nextjs
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/nextjs",USER="nextjs"

[program:falkordb]
command=sh -c './run.sh'
directory=%(ENV_FALKORDB_BIN_PATH)s
user=falkordb
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/falkordb",USER="falkordb"